# systemd-nvidia-entry/systemd-nvidia-entry.service
[Unit]
Description=Systemd Nvidia Entry

[Service]
Type=idle
ExecStart=/opt/MarechalLima/config-optimus.sh

[Install]
WantedBy=multi-user.target